# airtouch-webui.service
# systemd unit for the AirTouch web dashboard (Raspberry Pi OS).
# Assumes the repo is checked out at /home/frederico/pyairtouch-thermostat with a .venv inside —
# edit the paths and User= if yours differ, then see README.md for install steps.

[Unit]
Description=AirTouch web dashboard
After=network-online.target
Wants=network-online.target

[Service]
User=frederico
WorkingDirectory=/home/frederico/pyairtouch-thermostat
ExecStart=/home/frederico/pyairtouch-thermostat/.venv/bin/python /home/frederico/pyairtouch-thermostat/webui.py --config /home/frederico/pyairtouch-thermostat/config.toml
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
